# scripts/ecoplay-daily-summary.timer
[Unit]
Description=Run EcoPlay daily Lark summary every day

[Timer]
OnCalendar=*-*-* 18:00:00 Asia/Hong_Kong
Persistent=true
Unit=ecoplay-daily-summary.service

[Install]
WantedBy=timers.target
